# debian/postinst
#!/bin/sh
udevadm control --reload-rules
systemctl daemon-reload || true
systemctl enable rightclick.service || true
systemctl start rightclick.service || true

exit 0
